# image-builder/ubuntu/local-init.service
[Unit]
Description=Setup debug environment

[Service]
Type=oneshot
ExecStart=/root/local-init.sh

[Install]
WantedBy=multi-user.target